# etc/cron.daily/insight-cpanel-cron
#!/bin/bash
flock --exclusive \
  --nonblock /var/run/insight.running \
  --command /usr/local/lp/bin/insight-cpanel \
  &>/dev/null
