# usr/lib/systemd/system/mariadb@bootstrap.service.d/use_galera_new_cluster.conf
[Unit]
ConditionPathExists=
[Service]
Type=oneshot
Restart=no
ExecStartPre=
ExecStart=
ExecStart=/usr/bin/echo "Please use galera_new_cluster to start the mariadb service with --wsrep-new-cluster"
ExecStart=/usr/bin/false
ExecStartPost=