# etc/cron.d/csf_update
3 22 * * * root /usr/sbin/csf -u > /dev/null 2>&1
