# etc/cron.d/liquidweb_pkg_updates
16 14 * * * root yum clean all > /dev/null 2>&1; yum update -y > /dev/null 2>&1
