# etc/systemd/system/named.service.d/cpanel.conf
[Unit]
After=ipaliases.service
After=cpipv6.service
