# usr/lib/systemd/system/cryptsetup-pre.target
[Unit]
Description=Local Encrypted Volumes (Pre)
Documentation=man:systemd.special(7)
RefuseManualStart=yes
Before=cryptsetup.target