# usr/lib/systemd/system/cryptsetup.target
[Unit]
Description=Local Encrypted Volumes
Documentation=man:systemd.special(7)