# usr/lib/systemd/system/nss-lookup.target
[Unit]
Description=Host and Network Name Lookups
Documentation=man:systemd.special(7)
RefuseManualStart=yes